# systemd/openclaw-gateway-4.service
[Unit]
Description=OpenClaw Gateway-4
After=network.target

[Service]
Type=simple
User=yf
WorkingDirectory=/home/yf/.openclaw-4
ExecStart=/usr/local/bin/openclaw --profile 4 gateway
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
