# systemd/openclaw-gateway-7.service
[Unit]
Description=OpenClaw Gateway-7
After=network.target

[Service]
Type=simple
User=yf
WorkingDirectory=/home/yf/.openclaw-7
ExecStart=/usr/local/bin/openclaw --profile 7 gateway
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
