# systemd/openclaw-gateway-lawrence.service
[Unit]
Description=OpenClaw Gateway-lawrence
After=network.target

[Service]
Type=simple
User=yf
WorkingDirectory=/home/yf/.openclaw-lawrence
ExecStart=/usr/local/bin/openclaw --profile lawrence gateway
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
